# manufacture_auto/car/manufacture.service
[Unit]
Description=Manufacture

[Service]
ExecStart=/etc/systemd/system/manufacture.sh

[Install]
WantedBy=multi-user.target
