# manufacture_auto/angle/manufacture.service
[Unit]
Description=Angle arm

[Service]
ExecStart=/etc/systemd/system/manufacture.sh

[Install]
WantedBy=multi-user.target
