# manufacture_auto/delta/manufacture.service
[Unit]
Description=Delta manufacture

[Service]
ExecStart=/etc/systemd/system/manufacture.sh

[Install]
WantedBy=multi-user.target
